# root/etc/inittab
console::sysinit:/etc/rc.sysinit
console::askfirst:-/bin/sh
